# lib/supervisord.conf.d/0setup.conf
[program:setup]
command=/srv/bin/dbb --setup
autorestart=false
startsecs=0
startretries=0
user=docker
directory=/srv
