# lib/supervisord.conf.d/avahi-daemon.conf
[program:dbus]
command=dbus-daemon --system --nofork
directory=/
user=root

[program:avahi-daemon]
command=avahi-daemon
directory=/
user=root
