# lib/supervisord.conf.d/buildbot.conf
[program:buildmaster]
command=twistd --nodaemon --no_save -y buildbot.tac
directory=/srv/master
user=docker

[program:buildworker]
command=twistd --nodaemon --no_save -y buildbot.tac
directory=/srv/slave
user=docker
